# deployment/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid
user=appuser

[program:fastapi]
command=uvicorn backend.app:app --host 0.0.0.0 --port 8000 --workers 2 --log-level info
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/fastapi.log
stderr_logfile=/app/logs/fastapi_err.log

[program:streamlit]
command=python -m streamlit run streamlit_app.py
    --server.port=8501
    --server.address=0.0.0.0
    --server.headless=true
    --browser.gatherUsageStats=false
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/streamlit.log
stderr_logfile=/app/logs/streamlit_err.log
